# rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

export LD_LIBRARY_PATH=/home/pi/PythonProjects/Mjpg-Streamer/mjpg-streamer/mjpg-streamer-experimental/

/home/pi/PythonProjects/Mjpg-Streamer/mjpg-streamer/mjpg-streamer-experimental/mjpg_streamer -o "output_http.so -w /home/pi/PythonProjects/Mjpg-Streamer/mjpg-streamer/mjpg-streamer-experimental/www" -i "input_raspicam.so -x 640 -y 480 -fps 2 -quality 50" &

python3 /home/pi/PythonProjects/gui/app.py &

exit 0
